# deploy/eva-gateway.service
[Unit]
Description=Eva Gateway
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/eva
Environment="PATH=/opt/eva/.venv/bin"
EnvironmentFile=/opt/eva/.env
ExecStart=/opt/eva/.venv/bin/gunicorn -w 2 -b 127.0.0.1:18790 src.gateway:app
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
